# etc/systemd/system/ollama.service.d/override.conf
[Service]
CPUQuota=450%
Environment="OLLAMA_HOST=0.0.0.0"
