# etc/profile.d/gnupg.sh
export GNUPGHOME="$XDG_STATE_HOME"/gnupg
