# etc/profile.d/mariadb.sh
export MYSQL_HISTFILE="$XDG_STATE_HOME"/mariadb/.mariadb_history
