# etc/profile.d/postgresql.sh
export PSQL_HISTORY="$XDG_STATE_HOME"/psql/history
