# etc/profile.d/rustup.sh
append_path '/usr/lib/rustup/bin'
append_path $HOME/.cargo/bin
export PATH
