# etc/profile.d/ssh.sh
export SSH_HOME="$XDG_CONFIG_HOME"/ssh
export SSH_CONFIG="$XDG_CONFIG_HOME"/ssh/config
